# supervisord.conf
# Supervisor Configuration for Unified Container
# Use with Dockerfile.unified to manage multiple processes in one container
# Save as: supervisord.conf in project root

[supervisord]
nodaemon=true
logfile=/app/logs/supervisord.log
pidfile=/var/run/supervisord.pid

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
startsecs=10
stopwait=10
stdout_logfile=/app/logs/nginx.log
stderr_logfile=/app/logs/nginx_error.log

[program:gateway]
command=bash -c "cd /app/backend && PYTHONPATH=. uvicorn app.gateway.app:app --host 0.0.0.0 --port 8001 --workers 2"
directory=/app/backend
autostart=true
autorestart=true
startsecs=10
stopwait=30
stdout_logfile=/app/logs/gateway.log
stderr_logfile=/app/logs/gateway_error.log
environment=PYTHONUNBUFFERED=1,CI=true

[program:langgraph]
command=bash -c "cd /app/backend && uv run langgraph dev --no-browser --allow-blocking --host 0.0.0.0 --port 2024"
directory=/app/backend
autostart=true
autorestart=true
startsecs=15
stopwait=30
stdout_logfile=/app/logs/langgraph.log
stderr_logfile=/app/logs/langgraph_error.log
environment=PYTHONUNBUFFERED=1,CI=true

[program:frontend]
command=bash -c "cd /app/frontend && pnpm run start"
directory=/app/frontend
autostart=true
autorestart=true
startsecs=10
stopwait=30
stdout_logfile=/app/logs/frontend.log
stderr_logfile=/app/logs/frontend_error.log
environment=NODE_ENV=production

[group:deerflow]
programs=nginx,gateway,langgraph,frontend
priority=999
